# device-service/systemd/vehicle-overseer.service
[Unit]
Description=vehicle overseer client
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vehicle-overseer/app
EnvironmentFile=/etc/vehicle-overseer/device.env
ExecStart=/usr/bin/python3 -u /opt/vehicle-overseer/app/service.py run
Restart=always
RestartSec=3
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
