# updater/systemd/updater.service
[Unit]
Description=Vehicle Overseer updater (pull manifest + update)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/vehicle-overseer
EnvironmentFile=-/etc/vehicle-overseer/device.env
EnvironmentFile=-/etc/vehicle-overseer/updater.env
ExecStart=/usr/bin/python3 -u /opt/vehicle-overseer/updater.py
