# updater/systemd/updater.timer
[Unit]
Description=Vehicle Overseer updater timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=60min
RandomizedDelaySec=5min
Unit=vehicle-overseer-updater.service

[Install]
WantedBy=timers.target
